# socket/MyService.service
[Unit]
Description=Server

[Service]
Type =simple
ExecStart =/home/sofia/src/socket/server/build

[Install]
WantedBy=multi-user.target
